# configs/crontab
# m h  dom mon dow   command
*/5 * * * * cd /home/guardrailbb/buildbot; ./buildbot configs/guardrail
